# supervisord.d/icewm-session.ini
[program:icewm-session]
priority=10
autostart=true
user=%(ENV_WEBUSERNAME)s
command=bash -c 'sleep 3 && exec /usr/bin/icewm-session -d %(ENV_DISPLAY)s'
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
