# supervisord.d/sshd.ini
[program:sshd]
priority=0
autostart=true
command=/usr/sbin/sshd -D
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
