# supervisord.d/websockify.ini
[program:websockify]
priority=20
autostart=true
command=/usr/local/bin/websockify --ssl-only --cert /home/administrator/.ssh/selfpem --web=/usr/share/novnc/ 443 localhost:5901
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
