# supervisord.d/Xvnc.ini
[program:Xvnc]
priority=1
autostart=true
startsecs=2
user=%(ENV_WEBUSERNAME)s
command=Xvnc %(ENV_DISPLAY)s -geometry %(ENV_GEOMETRY)s -auth /home/%(ENV_WEBUSERNAME)s/.Xauthority -localhost -AlwaysShared -desktop cim-station -fp catalogue:/etc/X11/fontpath.d -pn -pn -rfbauth /home/%(ENV_WEBUSERNAME)s/.vnc/passwd -rfbport 5901
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
